# infra/vault/vps-overlay/etc/systemd/system/vault-to-ssm-sync.timer
[Unit]
Description=Vault → SSM sync timer — fires every 5 min
Documentation=https://github.com/hlebtkachenko/monorepo/blob/main/docs/plans/SECRETS-MIGRATION.md#m4--migrate-2-secrets-to-vault--aws-ssm-securestring-tunnel_token-direct-to-ssm

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30
Unit=vault-to-ssm-sync.service

[Install]
WantedBy=timers.target
